# infra/systemd/boats-com-scrape.timer
[Unit]
Description=OnlyHulls — daily boats.com scrape (residential egress via Elmo)
Documentation=https://github.com/overlordforthew/OnlyHulls

[Timer]
# Four fires per day, 6h apart, each with up to 30 minutes of jitter.
# Windows are chosen to catch inventory rotation across US, EU, and APAC business hours.
OnCalendar=*-*-* 03:00:00
OnCalendar=*-*-* 09:00:00
OnCalendar=*-*-* 15:00:00
OnCalendar=*-*-* 21:00:00
RandomizedDelaySec=1800
# If the machine was off at the scheduled time, catch up on next boot.
Persistent=true
Unit=boats-com-scrape.service

[Install]
WantedBy=timers.target
