# infra/systemd/boats-com-scrape.service
[Unit]
Description=OnlyHulls — scrape boats.com via Elmo and import to DB
Documentation=https://github.com/overlordforthew/OnlyHulls
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
# PAGES controls how many listing pages to request from boats.com. 15 = ~210 listings.
# Raise with care; keep it in the 5..30 range while boats_com is in source-policy "test".
Environment=PAGES=15
ExecStart=/root/projects/OnlyHulls/scripts/daily-scrape-boats-com.sh ${PAGES}
# Hard ceiling: boats.com pacing (30-75s × 15 pages) + import typically finishes in 15min.
# RuntimeMaxSec does not apply to Type=oneshot; TimeoutStartSec is the hang guard.
TimeoutStartSec=45min
KillMode=mixed
TimeoutStopSec=60s
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
# Hardening (kept narrow so docker exec / ssh / scp / notify.sh all still work).
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
# Best-effort logs land under /var/log/onlyhulls; journal also captures stdout/stderr.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=boats-com-scrape

[Install]
WantedBy=multi-user.target
